# packages/playout/etc/fk-playout.service
[Unit]
Description=Frikanalen Playout
After=syslog.target network.target

[Service]
Type=simple
User=fk-playout
WorkingDirectory=/opt/fk-playout
ExecStart=/opt/fk-playout/env/bin/python bin/integrated.py
Restart=on-abort

[Install]
WantedBy=multi-user.target
