# infra/roles/web/templates/fill_next_weeks_agenda.timer.j2
[Unit]
Description=Run fill_next_weeks_agenda daily (at 22)

[Timer]
OnCalendar=*-*-* 22:00:00
Persistent=true

[Install]
WantedBy=timers.target
